# bootstrap/conf/salt/state/consul/config/client/client-upstart.conf
description "Consul client process"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

respawn

setuid root
setgid root

exec consul agent -config-dir /etc/consul.d -ui-dir /etc/consul/ui